# cookbook/templates/default/upstart-github-connector-web.conf.erb
# THIS FILE IS MANAGED BY CHEF
# Local modifications will be discarded.

description "GitHub Connector"

start on runlevel [2]
stop on runlevel [016]

setuid <%= node['github_connector']['user'] %>
chdir <%= node['github_connector']['install_dir'] %>
exec <%= @rvm_path %>/bin/rvm <%= node['github_connector']['rvm_alias'] %> do puma -e production -b unix://<%= node['github_connector']['install_dir'] %>/tmp/sockets/puma.sock
reload signal SIGUSR2

respawn
